# scripts/firstrun
#!/bin/bash
# Firstrun setup for Raspberry Pi. Copied to /etc/init.d/firstrun
### BEGIN INIT INFO
# Provides:          firstrun
# Required-Start:    $networking
# Required-Stop:     $networking
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
### END INIT INFO

LIBNFCDIR=/home/pi/libnfc-1.7.1/
MAINDIR=/home/pi/src/vnw/main/
GODIR=/home/pi/go/
GOPATH=/home/pi/
SEMAPHORE=/var/opt/vnw-run
if [ ! -f $SEMAPHORE ]
then
  # Install our deps
	apt-get update
	apt-get install libusb-dev daemontools vim bzr golang
	# Install libnfc
	cd $LIBNFCDIR
	$LIBNFCDIR/configure
	make
	make install
	# Make a custom go
	cd $GODIR/src/
	./make.bash
	# Make and install our main program.
	cd $MAINDIR
	GOPATH=$GOPATH $GODIR/bin/go build main.go
	install main /bin/
  install scripts/lock /bin/
  install scripts/unlock /bin/
  
  # Scripts to make things run.
  echo "America/Los_Angeles" > /etc/timezone
  dpkg-reconfigure --frontend noninteractive tzdata
  # This is about the only non-idempotent thing in here.
  echo $'0\t0\t*\t*\t*\troot\t/bin/lock' >> /etc/crontab
	echo "export GOPATH=/home/pi/" >> /home/pi/.bashrc
  touch $SEMAPHORE
fi

# Start our program.
svscanboot &
exit 0
